# conf/servo-engine.upstart.conf

start on filesystem or runlevel [2345]
stop on shutdown

expect fork
respawn
script
	kore -c /usr/local/servo/conf/kore.conf
end script